# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
command=python3 -c "import uvicorn; uvicorn.run('main:app', host='0.0.0.0', port=8000, log_level='info')"
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend_err.log
environment=NODE_ENV="production"

[program:frontend]
command=node dist-server/server.mjs
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/frontend.log
stderr_logfile=/var/log/frontend_err.log
environment=NODE_ENV="production",PORT="5000",BACKEND_PORT="8000"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx_access.log
stderr_logfile=/var/log/nginx_error.log
